# installers/go-agent/sol/svc.xml
<?xml version="1.0"?>
<!-- *************************GO-LICENSE-START******************************
 * Copyright 2014 ThoughtWorks, Inc.
 *
 * Licensed under the Apache License, Version 2.0 (the "License");
 * you may not use this file except in compliance with the License.
 * You may obtain a copy of the License at
 *
 *    http://www.apache.org/licenses/LICENSE-2.0
 *
 * Unless required by applicable law or agreed to in writing, software
 * distributed under the License is distributed on an "AS IS" BASIS,
 * WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
 * See the License for the specific language governing permissions and
 * limitations under the License.
 *************************GO-LICENSE-END******************************* -->

<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="go-agent">

<service name="application/go/agent" type="service" version="1">
	<!--
	   Wait for network interfaces to be initialized.
	-->
        <dependency name="network" grouping="require_all" restart_on="none" type="service">
                <service_fmri value="svc:/milestone/network:default"/>
        </dependency> 

	<!--
	   Wait for all local filesystems to be mounted.
	-->
        <dependency name="filesystem-local" grouping="require_all" restart_on="none" type="service">
                <service_fmri value="svc:/system/filesystem/local:default"/>
        </dependency> 

        <exec_method type="method" name="start" exec="/lib/svc/method/go-agent start" timeout_seconds="60"/>
        <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60"/>
        <!--
        Don't restart server/agent if child coredumps
        http://blogs.everycity.co.uk/alasdair/2010/01/making-solaris-smf-ignore-core-dumps-in-child-processes/
        http://hub.opensolaris.org/bin/view/Community+Group+smf/faq
        -->
        <property_group name='startd' type='framework'>
                <propval name='ignore_error' type='astring'
                        value='core,signal' />
        </property_group>

        <instance name="default" enabled="true">
        	<method_context>
                	<method_credential user="go" group="daemon"/>
        	</method_context>

            <property_group name='go-agent' type='application'>
                <propval name='server' type='astring' value='127.0.0.1' /> 
            </property_group> 
        </instance>

        <stability value="Evolving"/>

        <template>
                <common_name>
                        <loctext xml:lang="C">
                                Go Agent
                        </loctext>
                </common_name>
                <documentation>
                        <doc_link name="thoughtworks.com/products/support" uri="http://www.thoughtworks.com/products/support"/>
                </documentation>
        </template>

</service>
</service_bundle>
